# docker/etc/supervisor/conf.d/supervisord-app.conf
[program:node]
command=node api/server.js
directory=/home/app/webapp/api
priority=3
stdout_logfile=syslog
stderr_logfile=syslog
